# tmux/log_tmux_panes.cron
*/10 * * * * $HOME/.dotfiles/tmux/log_tmux_panes.sh
